# configs/nono.service
[Unit]
Description=Nono, voice assistent command tool.
After=default.target

[Service]
WorkingDirectory=/home/natsu/Proyectos/onebyt/learnings/python/nono
ExecStart=/usr/bin/python3 main.py
Restart=always

[Install]
WantedBy=default.target
